# consul/root/etc/s6/consul/finish
#!/usr/bin/env bash

consul leave

# if Consul dies, this will force the entire container to stop
s6-svscanctl -t /etc/s6
