# redis/root/etc/s6/redis/run
#!/usr/bin/env bash

exec redis-server /etc/redis/redis.conf
